# deploy/wfm-scrape-pull.service
[Unit]
Description=Pull the CI-built wfm-scrape pipeline binary when scrape-latest changes
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/srv/wfm/pull-scrape.sh
