# deploy/wfm-scrape.service
[Unit]
Description=Refresh Warframe market.json from warframe.market + warframestat
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=wfm
Group=wfm
WorkingDirectory=/srv/wfm/app
ExecStart=/srv/wfm/run-scrape.sh
# Actual runtime crept past 3600s (item count + pacing grew) and systemd was
# killing the scrape a few hundred seconds short of writing market.json —
# every run for ~14h just re-served the last successful snapshot. Matches
# the GH Actions job's own timeout-minutes: 90 headroom.
TimeoutStartSec=5400
Nice=10

# Sandboxing — the scraper rewrites files under its repo, plus (phase 2b) the
# converter-shadow parity log. The '-' prefix keeps the unit startable if the
# log file is ever missing (setup-container.sh pre-creates it); the shadow write
# then simply no-ops rather than the whole service failing to start.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/wfm/app -/srv/wfm/shadow-parity.log
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallFilter=@system-service

# Ping a healthcheck on failure (optional — set the URL or delete this line).
# OnFailure=
